# obed@.service.in
[Unit]
Description=OBE C-@@ Encoder Master %i
After=obe-agent.service

[Service]
Type=simple
WorkingDirectory=/opt/obe-c@@
ExecStart=/opt/obe-c@@/bin/obed %i
Restart=always

[Install]
WantedBy=multi-user.target
