# obed1.conf
description "OBE D-100 Encoder Library"

start on started obed
stop on stopped obed

console log
respawn

exec start-stop-daemon --start --chdir /opt/obe-d100/ --exec /opt/obe-d100/bin/obed 1 --
