# debian/obed.upstart.in
description "OBE C-@@ Encoder Master"

start on filesystem or runlevel [2345]
stop on run level [!2345]

console log

pre-start script
    for i in `seq 0 7`; do
        start obed-worker id=$i || true
    done
end script
