# debian/obed.upstart
description "OBE D-100 Encoder Library"
 
start on filesystem or runlevel [2345]
stop on run level [!2345]

console log 
respawn

exec start-stop-daemon --start --chdir /opt/obe-d100/ --exec /opt/obe-d100/bin/obed 0 --
