# debian/obed.upstart
description "OBE C-100 Encoder Master"

start on filesystem or runlevel [2345]
stop on run level [!2345]

console log

pre-start script
    for i in `seq 0 7`; do
        start obed-worker id=$i cpu=$(($i & 1)) || true
    done
end script
